# service/sync_ipod@.service
[Unit]
Description=Syncing Ipod (/dev/%I)
After=systemd-udevd.service

[Service]
Type=notify
NotifyAccess=main
ExecStart=/usr/local/bin/sync_ipod /dev/%I
